# cmall.service.in
[Unit]
Description=cmall Core Daemon
After=network.target
Requires=postgresql.service

[Service]
ExecStartPre=/usr/bin/mkdir -p /var/log/cmall
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/cmall --db_host 127.0.0.1 --http 127.0.0.1:9797
WorkingDirectory=/var/log/cmall

LimitNOFILE=900000
Restart=always
RestartSec=5
MemoryMax=50G

[Install]
WantedBy=multi-user.target
